# deploy/weekly-board.service
[Unit]
Description=Weekly Board FastAPI Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/weekly-board
EnvironmentFile=/opt/weekly-board/.env
ExecStart=/opt/weekly-board/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
